# config/sensorReporter.service
[Unit]
Description=Reports presence status over MQTT and REST

[Service]
Type=simple
ExecStart=/opt/sensorReporter/sensorReporter.py /opt/sensorReporter/sensorReporter.ini
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

